# instances/pve-dev2/tree/etc/cron.d/trim-ssd

# trim our ssd daily
5 11 * * * root fstrim -v /mnt/ssd


